# src/systemd/nwg-quick.target
[Unit]
Description=NotWireGuard Tunnels via nwg-quick(8)
